# krb/templates/kadminhttp.service.erb
[Unit]
Description=kadminhttp service

[Service]
ExecStart=/usr/bin/python /opt/kadminhttp/kadminhttp.py
WorkingDirectory=/tmp
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=kadminhttp
Restart=always

[Install]
WantedBy=multi-user.target
