# rediser/templates/rediser.service.erb
[Unit]
Description=Rediser service

[Service]
Type=no
ExecStart=/usr/bin/python /opt/rediser/rediser7.py
WorkingDirectory=/tmp
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=rediser
Restart=always
User=<%= @service_user %>
Group=<%= @service_user %>

