# rediser/templates/rediser.service.erb
[Unit]
Description=Rediser service

[Service]
ExecStart=/usr/bin/python /opt/rediser/rediser_master.py
WorkingDirectory=/tmp
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=rediser
Restart=always
User=<%= @service_user %>
Group=<%= @service_user %>

[Install]
WantedBy=multi-user.target
