# rediser/templates/rediser.service.erb
[Unit]
Description=Rediser service

[Service]
Type=notify
ExecStart=/usr/sbin/rsyslogd -n -f /etc/rediser.conf -i /var/run/rediser/rediser.pid
StandardOutput=null
Restart=on-failure
User=<%= @service_user %>
Group=<%= @service_user %>

